# nexus.conf
[supervisord]
nodaemon=true

[program:nexus]
user=nexus
command=/bin/bash -c "/opt/nexus/bin/nexus console"
redirect_stderr=true
